# etc/services.d/cachet-monitor/run
#!/usr/bin/with-contenv sh

cachet-monitor -c /etc/monitoring/cachet-monitor.yaml